# overlay/airootfs/etc/xdg/autostart/win11-firstboot.desktop
[Desktop Entry]
Type=Application
Name=Win11Craft First Boot
Exec=/usr/local/bin/win11-firstboot
OnlyShowIn=KDE;
X-KDE-autostart-phase=1